# deploy/systemd/cdi-web.service
[Unit]
Description=Revenue Intelligence OS Web
After=network-online.target cdi-api.service
Wants=network-online.target

[Service]
Type=simple
User=cdi
Group=cdi
WorkingDirectory=/opt/revenue-intelligence-os
EnvironmentFile=/etc/cdi/revenue-intelligence-os.env
ExecStart=/opt/revenue-intelligence-os/scripts/production-start-web.sh
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
